# config/gunicorn.pyfindmespot.socket
[Unit]
Description=gunicorn.pyfindmespot.socket

[Socket]
ListenStream=/web/pyfindmespot/pyfindmespot.socket

[Install]
WantedBy=sockets.target
