# config/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/website/findmespot/app.socket

[Install]
WantedBy=sockets.target

